# bluetooth-server.service
[Unit]
Description=Bluetooth Audio Server
After=pulseaudio.service
Requires=pulseaudio.service
Requires=bluetooth-agent.service

[Service]
ExecStart=/root/bluetoothradio/connect.sh
Restart=always

[Install]
WantedBy=multi-user.target
Also=bluetooth-agent.service
