# systemd-logind.service
.include /lib/systemd/system/systemd-logind.service

[Service]
Environment=SYSTEMD_LOG_LEVEL=warning
